# templates/bitbucket.service.erb
[Unit]
Description=Atlassian Bitbucket Server Service
After=syslog.target network.target

[Service]
Type=forking
Environment="JAVA_HOME=<%= @javahome %>"
PIDFile=<%= @webappdir %>/work/catalina.pid
User=<%= @user %>
ExecStart=<%= @webappdir %>/bin/start-<%= @product %>.sh
ExecStop=<%= @webappdir %>/bin/stop-<%= @product %>.sh

[Install]
WantedBy=multi-user.target
